# nfs/netapp/netapp_diskscan.cron
#!/bin/sh
# $Id$
# $URL: $

# Change to 1 to enable diskscan
enableNetAppdiskscan=0

# Replace host with the NetApp storage appliance domain name
host=''

if [[ $enableNetAppdiskscan -eq 1 ]]; then
	/usr/local/sbin/netapp_diskscan.sh $host
fi
