# certify_check.cron
#!/bin/sh
# certify_check.cron
# $Id$
# $Date: $$

/usr/local/certify/check.py -l
