# certify.cron
# md5 hash check
10 23 * * * root /usr/local/certify/check.py -cl
# check
30 23 * * 0 root /usr/local/certify/check.py -l
# harden
50 23 1 * * root /usr/local/certify/harden.py -l


