# certify_harden.cron
#!/bin/sh
# $Id: certify_harden.cron 414 2014-04-12 14:11:11Z rohare $
# $URL: file:///usr/local/svn/certify/trunk/certify_harden.cron $

/usr/local/certify/harden.py -l
