# certify_harden.cron
#!/bin/sh
# certify_harden.cron
# $Id$
# $Date$

/usr/local/certify/harden.py -l
