# certify_check.cron
#!/bin/sh
# certify_check.cron

/usr/local/certify/check.py -l
