# certify_md5chk.cron
#!/bin/sh
# certify_md5chk.cron
# $Id$
# $Date: $

/usr/local/certify/check.py -cl
