# clamscan.cron
#!/bin/bash
enabled=1

if [ $enabled == "1" ]; then
	wget -r -l1 -np -nH --cut-dirs=3 --no-check-certificate "https://zdiv-yum/software/VendorSoftware/clam" -P /var/tmp/clamav -A "*.cvd"

	/usr/local/sbin/clamscan.sh $i &
fi
