# clamscan.cron
#!/bin/bash
enabled=1

if [ $enabled == "1" ]; then
	wget -r -l1 -np -nH --cut-dirs=3 --no-check-certificate "https://zdiv-yum/software/VendorSoftware/clam" -P /var/tmp/clamav -A "*.cvd"

	for i in /bin /boot /etc /home /lib /lib64 /opt /root /sbin /usr /var; do
		/usr/local/sbin/clamscan.sh $i &
	done
fi
