# clamscan.cron
#!/bin/bash
# clamscan.cron
#
enableClamav=0
clamavDataDir="/var/lib/clamav"

if [ ${clamavDataDir} != "/var/lib/clamav" ]; then
		wget -r -l1 -np -nH --cut-dirs=3 --no-check-certificate "https://zdiv-yum/software/VendorSoftware/clam" -P ${clamavDataDir} -A "*.cvd"

	/usr/local/sbin/clamscan.sh &
fi
