# diskcheck.cron
#!/bin/sh
# $Id: diskcheck.cron 505 2014-05-22 02:23:47Z rohare $
# $URL: file:///usr/local/svn/certify/trunk/diskcheck.cron $

# Change to 1 to enable diskcheck
enableDiskCheck=0

if [[ $enableDiskCheck -eq 1 ]]; then
	/usr/local/sbin/diskcheck.sh
fi
