# diskcheck.cron
#!/bin/sh
# diskcheck.cron
# $Id$
# $Date$

# Change to 1 to enable diskcheck
enableDiskCheck=0

if [[ $enableDiskCheck -eq 1 ]]; then
	/usr/local/sbin/diskcheck.sh
fi
