# nfs/ddn/ddn_diskscan.cron
#!/bin/sh
# $Id$
# $Date$

# Change to 1 to enable diskscan
enableDiskScan=0

# Replace host with the SUN storage appliance domain name

if [[ $enableDiskScan -eq 1 ]]; then
	/usr/local/sbin/toss_ddn_diskscan.sh
fi
